# roles/podman_mailserver/templates/docker-mailserver.container.j2
{{ ansible_managed | comment }}
# Quadlet container — docker-mailserver (SMTP/IMAP with TLS + fail2ban)
# Rootless Podman, pasta networking via pod.

[Container]
ContainerName={{ podman_mailserver_container_name }}
# renovate: datasource=docker depName=mailserver/docker-mailserver
Image={{ podman_mailserver_image }}:{{ podman_mailserver_version }}@{{ podman_mailserver_digest }}
Pod={{ podman_mailserver_pod_name }}.pod

AddCapability=NET_ADMIN
DropCapability=ALL

EnvironmentFile={{ podman_mailserver_home }}/config/mailserver.env

Volume={{ podman_mailserver_home }}/mail-data:/var/mail:Z
Volume={{ podman_mailserver_home }}/mail-state:/var/mail-state:Z
Volume={{ podman_mailserver_home }}/mail-log:/var/log/mail:Z
Volume={{ podman_mailserver_home }}/config:/tmp/docker-mailserver:Z
Volume={{ podman_mailserver_home }}/tls:/etc/dms/tls:ro,Z

HealthCmd=ss -tlnp | grep -q ':25'
HealthInterval=60s
HealthRetries=3
HealthStartPeriod=60s
HealthTimeout=5s

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=180
