# roles/notify_ip_change/templates/notify-ip-change.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Check for public IP changes every 5 minutes

[Timer]
OnCalendar=*-*-* *:00/5:00
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
