# roles/cloud_relay/templates/haproxy-podman.service.j2
{{ ansible_managed | comment }}
# Quadlet container unit — rootless Podman with pasta networking.
# Placed in ~/.config/containers/systemd/ for the {{ relay_haproxy_user }} user.
# See: podman-systemd.unit(5)
#
# Networking: pasta preserves source IPs (Podman 4.7+). HAProxy also sends
# PROXY protocol v2 at the application layer, so client IPs reach the backend
# regardless of the network driver.
#
# Port binding: privileged ports (25, 80, 443, etc.) are allowed by the kernel
# sysctl net.ipv4.ip_unprivileged_port_start (set via sysctl_configuration role),
# not by the NET_BIND_SERVICE capability (which is a user-namespace capability
# in rootless mode and has no effect on host port binding).

[Container]
ContainerName=haproxy-relay
Image={{ relay_haproxy_image }}
Network=pasta
ReadOnly=true
NoNewPrivileges=true

DropCapability=ALL

Volume={{ relay_haproxy_config_dir }}/haproxy.cfg:/usr/local/etc/haproxy/haproxy.cfg:ro,Z
Tmpfs=/run/haproxy

{% for fe in relay_haproxy_frontends %}
PublishPort={{ fe.bind_port }}:{{ fe.bind_port }}
{% endfor %}
{% if relay_haproxy_sni_routes | length > 0 %}
PublishPort={{ relay_haproxy_sni_port }}:{{ relay_haproxy_sni_port }}
{% endif %}
{% if relay_haproxy_http_redirect | default(true) and relay_haproxy_sni_routes | length > 0 %}
PublishPort=80:80
{% endif %}
{% if relay_haproxy_prometheus_enabled | default(false) %}
PublishPort={{ relay_haproxy_prometheus_port }}:{{ relay_haproxy_prometheus_port }}
{% endif %}

HealthCmd=haproxy -c -f /usr/local/etc/haproxy/haproxy.cfg
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=10s

[Service]
Restart=always
RestartSec=5
LogRateLimitIntervalSec=30s
LogRateLimitBurst=1000

[Install]
WantedBy=default.target
