# roles/kopia/templates/kopia-verify.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Weekly kopia snapshot verification

[Timer]
OnCalendar={{ kopia_verify_schedule }}
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
