# roles/kopia/templates/kopia-verify.service.j2
# {{ ansible_managed }}
[Unit]
Description=Kopia snapshot verification
After=network-online.target
Wants=network-online.target
OnFailure=kopia-failure-notify@%n.service

[Service]
Type=oneshot
WorkingDirectory=/root
EnvironmentFile=/etc/kopia.env
ExecStart={{ kopia_binary_path }} snapshot verify --verify-files-percent={{ kopia_verify_percent }}

# Hardening — same restrictions as kopia-snapshot
ReadOnlyPaths=/
ReadWritePaths=/root/.config/kopia /root/.cache/kopia /root/.ssh /tmp
NoNewPrivileges=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
