# roles/kopia/templates/kopia-snapshot.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Daily kopia backup snapshot

[Timer]
OnCalendar={{ kopia_snapshot_schedule }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
