# roles/kopia/templates/kopia-snapshot.service.j2
# {{ ansible_managed }}
[Unit]
Description=Kopia backup snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root
EnvironmentFile=/etc/kopia.env
ExecStart={{ kopia_binary_path }} snapshot create {{ kopia_snapshot_paths | join(' ') }}

# Hardening — Kopia runs as root (needs DAC read on /etc, /root, /var/lib/rancher)
ReadOnlyPaths=/
ReadWritePaths=/root/.config/kopia /root/.cache/kopia /root/.ssh /tmp
NoNewPrivileges=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
