# roles/hardening/templates/lynis-audit.service.j2
# {{ ansible_managed }}
[Unit]
Description=Lynis security audit
After=network-online.target
Wants=network-online.target
OnFailure=lynis-failure-notify@%n.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/lynis-audit.sh

# Hardening — ProtectSystem=full keeps /usr, /boot, /etc read-only; /var stays writable for sendmail
ProtectSystem=full
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
