# roles/mariadb-server/files/systemd/mysqld.service.d/syslog.conf
# according https://mariadb.com/kb/en/library/systemd/#configuring-mariadb-to-write-the-error-log-to-syslog
[Service]

StandardOutput=syslog
StandardError=syslog
SyslogFacility=daemon
SysLogLevel=err