# roles/docker-mailserver/files/scripts/update-resolver.conf
[program:updateresolverip]
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
command=/usr/local/bin/update-resolver-ip.sh
